# .config/systemd/user/beet.service
[Unit]
Description=beet
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/beet web

[Install]
WantedBy=multi-user.target
